# gcp/sawt.service
[Unit]
Description=Sawt WhatsApp ERP Gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=sawt
Group=sawt
WorkingDirectory=/opt/sawt
ExecStartPre=/opt/sawt/fetch-secrets.sh
ExecStart=/opt/sawt/sawt-gateway
EnvironmentFile=/opt/sawt/.env
Environment=GOMEMLIMIT=750MiB
Environment=GOGC=50

# --- resilience ---
Restart=always
RestartSec=5
StartLimitIntervalSec=0

# --- least privilege / sandboxing ---
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictSUIDSGID=true
LockPersonality=true
# writable paths the app / ExecStartPre actually need:
ReadWritePaths=/opt/sawt/voice-spool
# "-" tolerates the file not existing yet on a fresh install; fetch-secrets.sh
# still needs it pre-created once (touch + chown + chmod, see gcp/README.md)
# for the ProtectSystem=strict bind-mount below to attach to it at all.
ReadWritePaths=-/opt/sawt/.env

[Install]
WantedBy=multi-user.target
